# deploy/fase-0/systemd/heartbeat.service
# /etc/systemd/system/heartbeat.service
# heartbeat.py decide primero, en Python puro, si hay algo elegible en
# agent_inbox. Solo si lo hay invoca la sesión efímera claude --print que
# reclama y procesa los items del inbox — evita pagar una llamada al modelo
# en cada disparo del timer (cada 5 min) cuando el inbox está vacío.
[Unit]
Description=<Agent> heartbeat job
After=network-online.target postgresql.service

[Service]
Type=oneshot
User=<agent>
WorkingDirectory=/home/<agent>/claude
Environment=HOME=/home/<agent>/claude
Environment=PATH=/home/<agent>/apps/bin:/home/<agent>/claude/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=<AGENT>_CONTEXT=heartbeat
EnvironmentFile=/etc/<agent>/secrets.env
TimeoutStartSec=4min
# --strict-mcp-config es obligatorio: la sesión efímera NO debe heredar el plugin
# de Telegram (evita la colisión de polling con la sesión principal). heartbeat.py
# reproduce el mismo flag al invocar claude --print cuando hace falta.
ExecStart=/usr/bin/python3 /home/<agent>/workspace/scripts/lib/heartbeat.py
